# traffic-monitor.service
[Unit]
Description=Traffic Monitor Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/traffic-monitor
ExecStart=/usr/local/bin/traffic-monitor
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=traffic-monitor

# Environment
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# Security
NoNewPrivileges=false
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
